# Fall_2018/week2_playbook/supervisor.conf
[supervisord]
nodaemon=true

# [program:bro]
# command=/usr/bin/broctl start
# autostart=true
# autorestart=true

[program:openssh]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true