# Fall_2018/week1_dvwa/supervisord-openssh.conf
[program:openssh]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true